# roles/result/files/result.service
[Unit]
Description=Poll web application

[Service]
User=root
Environment
WorkingDirectory=/root/result/
ExecStart=npm start
Restart=always